# compose/cron/crontab
* * * * * root python /cron/soccer.py >> /var/log/cron.log 2>&1
